# etc/update-motd.d/00-header
#!/bin/sh

cat /var/run/motd_local-$(basename $0)
